# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisor/supervisord.pid
childlogdir=/run/supervisor

[program:gateway]
user=cloudron
directory=/app/code/backend
environment=PATH=/usr/local/bin:/usr/bin:/bin,HOME=/app/data/home,DEER_FLOW_HOME=/app/data/deer-flow,DEER_FLOW_CONFIG_PATH=/app/data/config.yaml,DEER_FLOW_EXTENSIONS_CONFIG_PATH=/app/data/extensions_config.json,DEER_FLOW_REPO_ROOT=/app/code,DEER_FLOW_HOST_BASE_DIR=/app/data/deer-flow,DEER_FLOW_HOST_SKILLS_PATH=/app/code/skills,DEER_FLOW_CHANNELS_LANGGRAPH_URL=http://127.0.0.1:8001,DEER_FLOW_CHANNELS_GATEWAY_URL=http://127.0.0.1:8001,PYTHONPATH=/app/code/backend,CI=true
command=/bin/bash -lc 'cd /app/code/backend && exec uv run --no-sync uvicorn app.gateway.app:app --host 0.0.0.0 --port 8001 --workers "${GATEWAY_WORKERS:-2}"'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=5
priority=10

[program:frontend]
user=cloudron
directory=/app/code/frontend
environment=PATH=/usr/local/bin:/usr/bin:/bin,HOME=/app/data/home,NODE_ENV=production,PORT=3000,HOSTNAME=127.0.0.1,BETTER_AUTH_SECRET=%(ENV_BETTER_AUTH_SECRET)s,BETTER_AUTH_BASE_URL=%(ENV_BETTER_AUTH_BASE_URL)s,DEER_FLOW_INTERNAL_GATEWAY_BASE_URL=http://127.0.0.1:8001,DEER_FLOW_INTERNAL_LANGGRAPH_BASE_URL=http://127.0.0.1:8001
command=/bin/bash -lc 'exec corepack pnpm start'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startsecs=8
priority=20

[program:nginx]
user=cloudron
command=/usr/sbin/nginx -c /run/nginx/nginx.conf -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=30
